# src/logger/cockpit-logger.service
[Unit]
Description=Cockpit Logger

[Service]
