# src/logger/cockpit-logger-janitor.timer
[Unit]
Description=Cockpit Logger Janitor

[Timer]
OnCalendar=hourly
