# src/logger/cockpit-logger-janitor.service.in
[Unit]
Description=Cockpit Logger Janitor

[Service]
Type=oneshot
ExecStart=@libexecdir@/cockpit-logger-janitor /var/log/cockpit/archives 2
